# conf/supervisord.d/djangoblog.ini
[program:djangoblog]
command = /data/app/scripts/djangoblog_start.sh
user = root
autostart = true
autorestart = false
redirect_stderr = true
stdout_logfile = /data/logs/blog.log
stderr_logfile = /data/logs/blog_error.log
